# .unitfiles/system/twmnd.service
[Unit]
Description=Twmnd Service
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/twmnd

